# docs/setup/asb/systemd/asb-watch-binary-change.path
[Path]
Unit=asb-restart.service
PathChanged=/var/external/asb/bin/asb

[Install]
WantedBy=multi-user.target
