# docs/setup/asb/systemd/asb-watch-binary-change.path
[Path]
Unit=asb-restart.service
PathChanged=/var/external/asb/xmr-btc-swap/target/release/asb

[Install]
WantedBy=multi-user.target
