# docs/setup/asb/systemd/pull-and-build-asb.timer
[Unit]
Description=triggers the service that pulls and builds asb

[Timer]
OnCalendar=16:00
Persistent=true

[Install]
WantedBy=timers.target
